# overlay/etc/supervisor/conf.d/nginx.conf
[program:nginx]
command=/opt/nginx/startnginx.sh
stdout_events_enabled=true
stderr_events_enabled=true
